# infra/cron/uptime_monitoring_report.cron

0 * * * * cd /opt/monitoring && sudo ./uptime-check-report.sh >> /var/log/uptime-check.log 2>&1